# crontab
* * * * * /usr/bin/sh /app/downloadVideos.sh