# crontab
* * * * * /bin/bash /app/downloadVideos.sh >> /var/log/download.log
* * * * * /bin/bash /app/transcodeVideos.sh >> /var/log/transcode.log